# prod/daemons/streaming/toshiSocket.conf
author "lauriermantel"
description "start and stop python bitcoin socket"
version "1.0"

start on started networking
stop on runlevel [!2345]

env APPUSER="root"
env APPDIR="/bin"
env APPBIN="toshiSocket"
env APPARGS=""

respawn

script
  exec su - $APPUSER -c "$APPDIR/$APPBIN $APPARGS"
end script
